# deploy/systemd/chat.service
[Unit]
Description=Rust WebSocket Chat
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/chat
EnvironmentFile=/opt/chat/.env
ExecStart=/opt/chat/target/release/Chat
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
